# common-file/gpiofan.service
[Unit]
Description=GPIO FAN

[Service]
Type=oneshot
ExecStart=/usr/sbin/gpiofan.sh

[Install]
WantedBy=multi-user.target
