# common-file/gpioled.service
[Unit]
Description=GPIO LED

[Service]
Type=oneshot
ExecStart=/usr/sbin/gpioled

[Install]
WantedBy=multi-user.target
